# web/root/etc/supervisor/conf.template.d/supervisord.app.conf
[program:app]
command=${APP_EXECUTABLE_PATH}
autostart=true
autorestart=true
