# web/root/etc/supervisor/conf.d.template/supervisord.app.conf
[program:app]
command=${APP_CMD}
autostart=true
autorestart=true
